# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/dichvulamsim.com/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/beobot.com/artisan schedule:run >> /dev/null 2>&1